# layer1-watchdog/com.openclaw.guardian.night.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.openclaw.guardian.night</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>{{HOME}}/.openclaw/scripts/openclaw-guardian/watchdog.py</string>
    </array>
    
    <!-- Night time: 00:00 - 07:00, once per hour -->
    <key>StartCalendarInterval</key>
    <array>
        <!-- Hour 0 -->
        <dict><key>Hour</key><integer>0</integer><key>Minute</key><integer>0</integer></dict>
        <!-- Hour 1 -->
        <dict><key>Hour</key><integer>1</integer><key>Minute</key><integer>0</integer></dict>
        <!-- Hour 2 -->
        <dict><key>Hour</key><integer>2</integer><key>Minute</key><integer>0</integer></dict>
        <!-- Hour 3 -->
        <dict><key>Hour</key><integer>3</integer><key>Minute</key><integer>0</integer></dict>
        <!-- Hour 4 -->
        <dict><key>Hour</key><integer>4</integer><key>Minute</key><integer>0</integer></dict>
        <!-- Hour 5 -->
        <dict><key>Hour</key><integer>5</integer><key>Minute</key><integer>0</integer></dict>
        <!-- Hour 6 -->
        <dict><key>Hour</key><integer>6</integer><key>Minute</key><integer>0</integer></dict>
    </array>
    
    <key>RunAtLoad</key>
    <false/>
    
    <key>StandardOutPath</key>
    <string>{{HOME}}/.openclaw/guardian/watchdog-launchd.log</string>
    
    <key>StandardErrorPath</key>
    <string>{{HOME}}/.openclaw/guardian/watchdog-launchd-error.log</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
